# src/rc.d/rc.local
#!/bin/bash
sh /etc/init.d/sethostname2.sh
sh /etc/init.d/userdataExecutor
